# jeffrey.service
[Unit]
Description=Jeffrey Personal Assistant Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fsalazar
WorkingDirectory=/home/fsalazar/jeffrey
ExecStart=/home/fsalazar/jeffrey/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Variables d'environnement chargées depuis le fichier .env
EnvironmentFile=/home/fsalazar/jeffrey/.env

[Install]
WantedBy=multi-user.target
